# docker/supervisord.conf
; ==============================================================================
; Supervisor Configuration for Laravel Application
; ==============================================================================
; Orchestrates PHP-FPM, Nginx, and essential bootstrap tasks in production.
; ==============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ------------------------------------------------------------------------------
; Runtime Bootstrapping (Runs once on start)
; ------------------------------------------------------------------------------
[program:laravel-init]
command=/bin/sh -c "/usr/local/bin/php /var/www/html/artisan storage:link --ansi --force || true; /usr/local/bin/php /var/www/html/artisan config:cache --ansi --quiet || true; /usr/local/bin/php /var/www/html/artisan route:cache --ansi --quiet || true; /usr/local/bin/php /var/www/html/artisan view:cache --ansi --quiet || true"
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

; ------------------------------------------------------------------------------
; PHP-FPM Process
; ------------------------------------------------------------------------------
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ------------------------------------------------------------------------------
; Nginx Process
; ------------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ------------------------------------------------------------------------------
; Optional: Laravel Queue Worker (Uncomment to enable queue processing)
; ------------------------------------------------------------------------------
;[program:laravel-worker]
;command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
;autostart=true
;autorestart=true
;numprocs=1
;redirect_stderr=true
;stdout_logfile=/var/log/supervisor/laravel-worker.log

; ------------------------------------------------------------------------------
; Optional: Laravel Scheduler (Uncomment to enable task scheduler cron-like loop)
; ------------------------------------------------------------------------------
;[program:laravel-scheduler]
;command=/bin/sh -c "while true; do /usr/local/bin/php /var/www/html/artisan schedule:run --verbose --no-interaction & sleep 60; done"
;autostart=true
;autorestart=true
;redirect_stderr=true
;stdout_logfile=/var/log/supervisor/laravel-scheduler.log
